# launchd/com.atlas.podcast-prompt-832weekends.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.atlas.podcast-prompt-832weekends</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/envchain</string>
		<string>atlas</string>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/printer/atlas/tools/podcast/weekly_prompt.py</string>
		<string>--podcast</string>
		<string>832weekends</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/printer/atlas</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>5</integer>
		<key>Hour</key>
		<integer>17</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/printer/atlas/logs/podcast_prompt_832weekends.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/printer/atlas/logs/podcast_prompt_832weekends.log</string>
</dict>
</plist>
