# launchd/com.atlas.rotary-print-agenda.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.atlas.rotary-print-agenda</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/envchain</string>
    <string>atlas</string>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/printer/atlas/tools/rotary/print_agenda.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/printer/atlas</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/printer/atlas/logs/rotary_print_agenda.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/printer/atlas/logs/rotary_print_agenda.log</string>
</dict>
</plist>
