# launchd/com.atlas.podcast-script-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.atlas.podcast-script-poller</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/envchain</string>
		<string>atlas</string>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/printer/atlas/tools/podcast/script_approver.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/printer/atlas</string>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>StandardOutPath</key>
	<string>/Users/printer/atlas/logs/podcast_script_poller.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/printer/atlas/logs/podcast_script_poller.log</string>
</dict>
</plist>
