# launchd/com.atlas.wa-dui-bill-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.atlas.wa-dui-bill-tracker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/envchain</string>
    <string>atlas</string>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/printer/atlas/tools/legal/wa_dui_bill_tracker.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/printer/atlas</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/printer/atlas/logs/wa_dui_bill_tracker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/printer/atlas/logs/wa_dui_bill_tracker.log</string>
</dict>
</plist>
